# nutries/nutries.ini
[program:nutries]
command=/apps/nutries/virtenv/bin/gunicorn --bind 127.0.0.1:8000 --workers 3 foodtasker.wsgi
directory=/apps/nutries/
autostart=true
autorestart=true
