# nutries/postgresql.ini
[program:postgresql]
user=postgres
command=/usr/pgsql-10/bin/postgres -D /var/lib/pgsql/10/data/ -c config_file=/var/lib/pgsql/10/data/postgresql.conf
autostart=true
autorestart=true
